# ansible/templates/ddns.service
[Unit]
Description=Dynamically update DNS host records.
After=network.target

[Service]
User={{ usr_name }}
Group={{ usr_name }}
ExecStart=/home/{{ usr_name }}/porkbun-ddns-pydantic/.venv/bin/python /home/{{ usr_name }}/porkbun-ddns-pydantic/main.py
WorkingDirectory=/home/{{ usr_name }}/porkbun-ddns-pydantic/
Restart=always
RestartSec=10
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target
